# systemd/testrr-storage-prune.timer
[Unit]
Description=Run testrr storage prune nightly

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
