# systemd/testrr.service
[Unit]
Description=testrr — test report hub
Documentation=https://github.com/pipelines-as-code/testrr
After=network-online.target
Wants=network-online.target
RequiresMountsFor=%t/containers

[Service]
Type=simple
Restart=on-failure
RestartSec=5s

# Runtime directory for the data volume — adjust to taste
Environment=TESTRR_DATA_DIR=/var/lib/testrr/data
Environment=TESTRR_ADDR=:8080

# Pull image from GHCR — set your image here
Environment=TESTRR_IMAGE=ghcr.io/your-org/testrr:latest

ExecStartPre=-/usr/bin/podman rm -f testrr
ExecStart=/usr/bin/podman run \
    --name testrr \
    --rm \
    --pull=newer \
    \
    --network=host \
    \
    --volume ${TESTRR_DATA_DIR}:/home/testrr/data:z,rw \
    \
    --env TESTRR_ADDR=${TESTRR_ADDR} \
    --env TESTRR_DATA_DIR=/home/testrr/data \
    \
    --user 1000:1000 \
    --userns=keep-id \
    \
    --cap-drop=ALL \
    --security-opt=no-new-privileges \
    --security-opt=label=type:container_runtime_t \
    --read-only \
    --tmpfs /tmp:rw,noexec,nosuid,size=64m \
    \
    ${TESTRR_IMAGE}

ExecStop=/usr/bin/podman stop testrr

# --- systemd hardening ---
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectClock=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadWritePaths=/var/lib/testrr/data
ReadOnlyPaths=/

[Install]
WantedBy=multi-user.target
